# config/s6/nginx/run
#!/bin/sh

# 启动Nginx
exec nginx -g "daemon off;"
